# files/etc/s6-overlay/s6-rc.d/zebra-server/run
#!/command/execlineb -P
with-contenv

importas -D "" elastic USE_ELASTICSEARCH

if -n { test "${elastic}" = "true" }

fdmove -c 2 1

s6-envdir /etc/koha-envvars

s6-setuidgid default-koha

/usr/bin/zebrasrv -v 1 -f /etc/koha/sites/default/koha-conf.xml

