# files/etc/s6-overlay/s6-rc.d/worker/run
#!/command/execlineb -P

fdmove -c 2 1

s6-envdir /etc/koha-envvars

s6-setuidgid default-koha

/usr/bin/perl /usr/share/koha/bin/background_jobs_worker.pl 
