# files/etc/s6-overlay/s6-rc.d/plack/run
#!/command/execlineb -P
fdmove -c 2 1

s6-envdir /etc/koha-envvars

/usr/bin/starman -M FindBin --socket /var/run/koha/default/plack.sock --max-requests 50 --workers 2 -E deployment --user=default-koha --group=default-koha --access-log /var/log/koha/default/plack.log --error-log /var/log/koha/default/plack-error.log "/etc/koha/plack.psgi" 
